# src/Services/ChristmasPi.service
[Unit]
Description=ChristmasPi
After=network.target

[Service]
ExecStart=

[Install]
WantedBy=multi-user.target